# skills/faion/knowledge/infra/cicd-tls-renewal-automation/templates/certbot-renew.service
# purpose: systemd service invoking certbot renew + deploy hook
# consumes: env + secrets-manager refs
# produces: systemd unit
# depends-on: content/01-core-rules.xml
# token-budget-impact: ~150 tokens when loaded as context

[Unit]
Description=cicd-tls-renewal-automation unit skeleton

[Install]
WantedBy=multi-user.target
