# skills/faion/knowledge/solo/infra/server-craft/systemd-user-services/templates/target.service
# ~/.config/systemd/user/nero.target
# Target unit: group NERO services so they start/stop together
#
# Usage:
#   systemctl --user enable nero.target    (sets as boot target)
#   systemctl --user start nero.target     (starts all Wants= services)
#   systemctl --user stop nero.target      (stops all Wants= services)
#   systemctl --user restart nero.target   (restarts all)

[Unit]
Description=NERO Platform
# List all services that belong to this platform
Wants=nero-core.service nero-channel-web.service nero-channel-tg.service nero-web.service

[Install]
WantedBy=default.target
