# skills/faion/knowledge/solo/infra/server-craft/systemd-user-services/templates/telegram-bot.service
[Unit]
Description=NERO Channel Telegram (aiogram bot)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=/srv/nero/nero-channel-tg/src
EnvironmentFile=/home/nero/workspace/.env

ExecStart=/srv/nero/nero-channel-tg/.venv/bin/python -m nero_channel_tg.main

Restart=on-failure
RestartSec=5

# Lightweight bot: limited resources
MemoryMax=1G
MemoryHigh=768M
MemorySwapMax=128M
CPUQuota=100%
LimitNOFILE=4096
TasksMax=64

OOMScoreAdjust=-200
OOMPolicy=stop

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
