# skills/faion/knowledge/solo/infra/server-craft/systemd-user-services/templates/timer-pair.service
# purpose: Template fixture for systemd-user-services: timer-pair.service
# consumes: content/01-core-rules.xml
# produces: structured artefact
# depends-on: content/02-output-contract.xml
# token-budget-impact: small
# Timer unit: schedule the service below
# ~/.config/systemd/user/backup.timer
[Unit]
Description=Daily backup timer

[Timer]
# Run once daily at midnight
OnCalendar=daily

# Run missed execution once after boot (e.g., server was down at midnight)
Persistent=true

# Spread load: random delay up to 5 minutes
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

# ─────────────────────────────────────────────────────────────
# Service unit (save as separate file: backup.service)
# ~/.config/systemd/user/backup.service
# ─────────────────────────────────────────────────────────────
# [Unit]
# Description=Daily backup
#
# [Service]
# Type=oneshot
# WorkingDirectory=/home/nero/workspace
# EnvironmentFile=/home/nero/workspace/.env
# ExecStart=/home/nero/workspace/scripts/backup.sh
# StandardOutput=journal
# StandardError=journal
#
# # No [Install] section needed for timer-activated services

# ─────────────────────────────────────────────────────────────
# Management commands:
#   systemctl --user enable --now backup.timer
#   systemctl --user list-timers
#   systemctl --user status backup.timer
#   journalctl --user -u backup.service    # view run history
