# skills/faion/knowledge/pro/infra/cicd-engineer/backup-filesystem-restic/templates/backup.service
# purpose: systemd service skeleton invoking backup-restic.sh
# consumes: env + secrets-manager refs
# produces: systemd unit
# depends-on: content/01-core-rules.xml
# token-budget-impact: ~150 tokens when loaded as context

[Unit]
Description=backup-filesystem-restic unit skeleton

[Install]
WantedBy=multi-user.target
