# skills/faion/knowledge/pro/dev/software-developer/php-laravel-queues/templates/supervisor.ini
; purpose: Supervisor config for Horizon-managed queue workers per named queue
; consumes: see content/02-output-contract.xml inputs
; produces: artefact conforming to content/02-output-contract.xml
; depends-on: content/01-core-rules.xml
; token-budget-impact: ~250 tokens when loaded as context

[program:laravel-horizon]
process_name=%(program_name)s
command=php /var/www/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/horizon.log
stopwaitsecs=3600

[program:laravel-worker-mail]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work redis --queue=mail --sleep=3 --tries=5 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/worker-mail.log

[program:laravel-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/worker-default.log
