# skills/faion/knowledge/pro/dev/software-developer/php-laravel-queues/templates/supervisor.ini
; /etc/supervisor/conf.d/laravel-worker.conf
; Queue worker daemon with priority queues and graceful restart.
; Run: supervisorctl reread && supervisorctl update && supervisorctl start laravel-worker:*

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work redis --queue=high,default,notifications,low --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/laravel-worker.log
; Allow 3600s for in-flight jobs to finish before force-kill.
stopwaitsecs=3600
