# skills/faion/knowledge/solo/infra/server-craft/health-checks-autoheal/templates/watcher.service
# ~/.config/systemd/user/watcher.service
# systemd user service for the auto-heal watcher process
# Enable: systemctl --user enable watcher.service
# Start:  systemctl --user start watcher.service

[Unit]
Description=Auto-Heal Watcher
After=network.target

[Service]
Type=simple
ExecStart=/srv/nero/watcher/.venv/bin/python3 /srv/nero/watcher/watcher.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
