# skills/faion/knowledge/pro/dev/backend-enterprise/php-laravel-queues/templates/laravel-worker.service
[Unit]
Description=Laravel queue worker (%i)
After=network.target redis.service

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=3

# --max-time bounds memory leaks; --max-jobs recycles workers after N jobs
ExecStart=/usr/bin/php /var/www/app/artisan queue:work redis \
  --queue=%i \
  --sleep=3 \
  --tries=3 \
  --max-time=3600 \
  --max-jobs=1000

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Enable named queue workers:
# systemctl enable --now laravel-worker@orders
# systemctl enable --now laravel-worker@notifications
# systemctl enable --now laravel-worker@low
