# skills/faion/knowledge/pro/dev/software-developer/ruby-sidekiq-jobs/templates/sidekiq.service
# /etc/systemd/system/sidekiq.service
# Systemd unit for daemonized Sidekiq with graceful shutdown (SIGTERM).
# Reload: systemctl daemon-reload && systemctl restart sidekiq

[Unit]
Description=Sidekiq Background Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/app/current
ExecStart=/usr/local/bin/bundle exec sidekiq -e production -C config/sidekiq.yml
Restart=always
KillSignal=SIGTERM
# Allow 30s for in-flight jobs to finish before force-kill.
TimeoutStopSec=30
User=deploy
Environment=RAILS_ENV=production

[Install]
WantedBy=multi-user.target
