# skills/faion/knowledge/solo/dev/software-developer/django-celery/templates/celery-worker.service
# purpose: systemd unit for Celery worker
# consumes: See content/02-output-contract.xml inputs
# produces: artefact conforming to content/02-output-contract.xml
# depends-on: content/01-core-rules.xml
# token-budget-impact: ~200-1000 tokens when loaded as context
# /etc/systemd/system/celery-emails.service
# Systemd unit for a single-queue Celery worker.
# Adjust: User, WorkingDirectory, ExecStart queue name, concurrency.

[Unit]
Description=Celery worker (emails queue)
After=network.target redis.service

[Service]
Type=simple
User=app
Group=app
EnvironmentFile=/etc/app/celery.env
WorkingDirectory=/srv/app
ExecStart=/srv/app/.venv/bin/celery -A config worker \
  -Q emails \
  -n emails@%%h \
  --concurrency=4 \
  --max-tasks-per-child=200 \
  --soft-time-limit=300 \
  --time-limit=360 \
  --without-gossip \
  --without-mingle
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
