# skills/faion/knowledge/backend/ruby-sidekiq-jobs/templates/sidekiq.service
# purpose: Legacy template for the ruby-sidekiq-jobs methodology.
# consumes: inputs declared in ruby-sidekiq-jobs/AGENTS.md prerequisites.
# produces: working code/config aligned with content/01-core-rules.xml.
# depends-on: content/02-output-contract.xml schema for output shape.
# token-budget-impact: ~600 tokens when loaded as reference.
[Unit]
Description=sidekiq
After=syslog.target network.target

[Service]
Type=notify
WatchdogSec=10
User=deploy
WorkingDirectory=/var/www/app/current
ExecStart=/usr/local/bin/bundle exec sidekiq -e production -C config/sidekiq.yml
ExecReload=/bin/kill -TSTP $MAINPID
Restart=always
RestartSec=1
KillMode=mixed
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
